# images/fedora-systemd/services/flyssh/flyssh.service
[Unit]
Description=Fly SSH Service
After=network.target

[Service]
# Environment variables are provided via Fly secrets
ExecStart=/usr/local/bin/flyssh server -debug -dev -port 8080
# Add more verbose logging
Environment=RUST_LOG=debug
# Ensure proper restart behavior
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target 