# images/ubuntu-s6/services/sshd/run
#!/bin/sh

# Start sshd in the foreground
exec /usr/sbin/sshd -D 