# images/ubuntu-s6/services/container-alive/run
#!/bin/sh

# Keep container alive
while true; do
    sleep 3600
done 