# images/fedora-systemd/services/show-status.service
[Unit]
Description=Show systemd service status after boot
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'systemctl list-units --type=service --all && echo "=== Service Tree ===" && systemctl list-dependencies --all'
StandardOutput=inherit
StandardError=inherit
LogLevelMax=debug

[Install]
WantedBy=multi-user.target 